# .config/dev/supervisor_app.conf
[program:uwsgi]
command=pipenv run uwsgi --ini /srv/project/.config/dev/uwsgi.ini

[program:nginx]
command=nginx