# .config/dev/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /srv/Brunch-RSSFeed/.config/dev/uwsgi.ini

[program:nginx]
command=nginx